# code/homebridge.timer
[Unit]
Description=Delay

[Timer]
OnBootSec=60


[Install]
WantedBy=multi-user.target
